# usb-auto/90-k2h-usb.rules
ACTION=="add|remove", KERNEL=="sdb", SUBSYSTEMS=="usb", DRIVERS=="usb-storage", RUN+="/usr/local/bin/ctlbiotropyd"
