# usb-auto/90-k2h-usb.rules
ACTION=="add", KERNEL=="sdb", SUBSYSTEMS=="usb", DRIVERS=="usb-storage", MODE="0777", SYMLINK="k2myheart", RUN+="/usr/local/bin/ctlk2husb"
